# extra/gbr-commission.timer
[Unit]
Description=Recommission nodes

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
